# src/supervisord.conf
[supervisord]
nodaemon=true

[program:run-mysite.co.uk]
directory=/src/
command=node /src/index.js
autostart=true
autorestart=true
stderr_logfile=/var/log/run-mysite.co.uk.err.log